# misc/.config/systemd/user/battery-warning.service
[Unit]
Description=Battery warning check

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/battery_warning.sh

